# root/etc/s6-overlay/s6-rc.d/init-xteve-config/run
#!/usr/bin/with-contenv bash

# create directories
mkdir -p \
    /config/{cache,conf,backup,cron} \
    /tmp/{xteve}

# permissions
lsiown -R abc:abc \
    /tmp/xteve \
    /config \
    /config/*
